# deploy/ella-chat-publish.service
[Unit]
Description=Ella Chat Publish Bot
After=network.target

[Service]
Type=simple
User=swiri021
WorkingDirectory=/home/swiri021/ella-chat-publish
ExecStart=/home/swiri021/ella-chat-publish/venv/bin/python -m src.bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
