# deploy/llama-server.service
[Unit]
Description=llama-cpp-python LLM Server
After=network.target

[Service]
Type=simple
User=swiri021
WorkingDirectory=/home/swiri021/unsloth
ExecStart=/home/swiri021/unsloth/venv/bin/python -m llama_cpp.server --model /home/swiri021/unsloth/models/gemma-4-31B-it-uncensored-Q4_K_M.gguf --host 0.0.0.0 --port 8080 --n_ctx 8192 --n_gpu_layers -1 --flash_attn true
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
